# src/bhunter.service
[Unit]
Description = XXX

[Service]
ExecStart=/usr/bin/bhunter start
Type=fork

[Install]
WantedBy=multi-user.target
